# deployment/server/systemd/visiondrive-license.service
[Unit]
Description=VisionDrive License Plate GPU Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/visiondrive/app/algorithm
EnvironmentFile=%h/visiondrive/config/runtime.env
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
Environment=OPENBLAS_NUM_THREADS=2
ExecStart=%h/miniconda3/envs/visiondrive/bin/python -m license.main
Restart=on-failure
RestartSec=3
TimeoutStartSec=180

[Install]
WantedBy=visiondrive.target
