# deployment/server/systemd/visiondrive-vehicle.service
[Unit]
Description=VisionDrive PP-Vehicle GPU Tracking Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/visiondrive/app/algorithm
EnvironmentFile=%h/visiondrive/config/runtime.env
Environment=PYTHONUNBUFFERED=1
Environment=MPLCONFIGDIR=/tmp/visiondrive-matplotlib
ExecStart=%h/miniconda3/envs/visiondrive-vehicle/bin/python -m vehicle.main
Restart=on-failure
RestartSec=3
TimeoutStartSec=240

[Install]
WantedBy=visiondrive.target
