# deployment/server/systemd/visiondrive-police.service
[Unit]
Description=VisionDrive Police Gesture GPU Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/visiondrive/app/algorithm
EnvironmentFile=%h/visiondrive/config/runtime.env
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/miniconda3/envs/visiondrive/bin/python -m police.main
Restart=on-failure
RestartSec=3
TimeoutStartSec=180

[Install]
WantedBy=visiondrive.target

