# deployment/server/systemd/visiondrive.target
[Unit]
Description=VisionDrive Server Services
Wants=visiondrive-license.service visiondrive-vehicle.service visiondrive-police.service visiondrive-gesture.service visiondrive-backend.service visiondrive-jpeg.service visiondrive-frontend.service
After=network-online.target

[Install]
WantedBy=default.target
